# debian/lxd-tcp.socket
[Socket]
ListenStream=[::]:8443
Service=lxd.service

[Install]
WantedBy=sockets.target
